# api/.docker/prod/schedule/cron.d/crontab
* * * * * /usr/local/bin/php /var/www/html/artisan schedule:run >> /proc/1/fd/1 2>/proc/1/fd/2
