# docker/dev/supervisor/conf.d/cron.conf
[program:cron]
command=cron && tail -f /var/log/cron.log
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/cron.log
